# setupfiles/lpfm.service
# LPFM station systemd service file
#
# Link and enable with:
# % sudo ln -sf /home/pi/lpfm/setupfiles/lpfm.service /etc/systemd/system/lpfm.service
# % sudo systemctl daemon-reload
# % sudo systemctl enable lpfm.service
#
# Start:
# % sudo systemctl start lpfm.service
#
# Stop:
# % sudo systemctl stop lpfm.service
#
# Restart:
# % sudo systemctl restart lpfm.service
#
# View journal:
# % journalctl -u lpfm.service -f

[Unit]
Description=LPFM low-power FM station
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/lpfm
ExecStart=/home/pi/lpfm/venv/bin/python main.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
